# build/COPY_ROOT/etc/supervisor/supervisord/conf.d/comfyui_rp_api.conf
[program:comfyui_rp_api]
command=/opt/ai-dock/bin/supervisor-comfyui-rp-api.sh
process_name=%(program_name)s
numprocs=1
directory=/opt/serverless/providers/runpod
priority=100
autostart=true
startsecs=5
startretries=3
autorestart=unexpected
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/comfyui-rp-api.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=1
stderr_logfile=/dev/null
stderr_logfile_maxbytes=0
stderr_logfile_backups=0
environment=PROC_NAME="%(program_name)s"