# build/COPY_ROOT/etc/supervisor/supervisord/conf.d/comfyui.conf
[program:comfyui]
user=$USER_NAME
environment=PROC_NAME="%(program_name)s",USER=$USER_NAME,HOME=/home/$USER_NAME
command=/opt/ai-dock/bin/supervisor-comfyui.sh
process_name=%(program_name)s
numprocs=1
directory=/opt/ComfyUI
priority=100
autostart=true
startsecs=5
startretries=3
autorestart=unexpected
stopsignal=TERM
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/supervisor/comfyui.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=1
stderr_logfile=/dev/null
stderr_logfile_maxbytes=0
stderr_logfile_backups=0
redirect_stderr=true